# alexa.conf
description "NodeJS Server"
author "Jeff"

env PORT=4433

start on runlevel [2345]
stop on runlevel [016]
respawn

setuid admin
chdir /volume1/homes/admin/echo/
exec node server.js